# upstart/nomad.conf
description "Nomad by HashiCorp"
start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
umask 022
chroot /
chdir /
limit core 0 0
#limit msgqueue <softlimit> <hardlimit>
#limit nice <softlimit> <hardlimit>
limit nofile 262144 262144
limit nproc 1048576 1048576
#limit rtprio <softlimit> <hardlimit>
#limit sigpending <softlimit> <hardlimit>
setuid root
setgid root
respawn
post-stop exec sleep 10

script
  # When loading default and sysconfig files, we use `set -a` to make
  # all variables automatically into environment variables.
  set -a
  [ -r /etc/default/nomad ] && . /etc/default/nomad
  set +a
  exec /usr/sbin/nomad "agent" "-config=/etc/nomad/"
end script
